# homebox-mcp/rootfs/etc/s6-overlay/s6-rc.d/homebox-mcp/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Homebox MCP Server
# Handles service shutdown
# ==============================================================================

bashio::log.info "Homebox MCP Server stopped"
