# homebox-mcp/rootfs/etc/s6-overlay/s6-rc.d/homebox-mcp/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Homebox MCP Server
# Runs the MCP server for Homebox inventory management
# ==============================================================================

bashio::log.info "Starting Homebox MCP Server..."

# Export configuration as environment variables
export HOMEBOX_URL="$(bashio::config 'homebox_url')"
export HOMEBOX_TOKEN="$(bashio::config 'homebox_token')"
export LOG_LEVEL="$(bashio::config 'log_level')"

bashio::log.info "Connecting to Homebox at: ${HOMEBOX_URL}"

# Run the MCP server
exec python3 /app/server.py
